# infra/systemd/home-sensors-collector-outside.service
[Unit]
Description=Home Sensors collector — Outside (192.168.1.36)
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/home-sensors/collectors
EnvironmentFile=/home/pi/home-sensors/.env
ExecStart=/usr/bin/python3 collect_data.py 192.168.1.36 Outside
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/HomeData/Logs/collect_data.log
StandardError=append:/home/pi/HomeData/Logs/collect_data.log

[Install]
WantedBy=multi-user.target
